# Yoyoo/project/backend/deploy/openclaw-http-bridge.service
[Unit]
Description=OpenClaw HTTP Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/yoyoo/backend
EnvironmentFile=-/root/yoyoo/backend/.env
ExecStart=/root/yoyoo/backend/.venv/bin/uvicorn app.services.openclaw_http_bridge:app --host 127.0.0.1 --port 18080
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
