# Yoyoo/project/backend/deploy/yoyoo-dingtalk-forwarder.service
[Unit]
Description=Yoyoo DingTalk Stream Forwarder
After=network.target yoyoo-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/yoyoo/backend
EnvironmentFile=-/root/yoyoo/backend/.env
Environment=YOYOO_DINGTALK_EVENTS_URL=http://127.0.0.1:8000/api/v1/dingtalk/events
Environment=NODE_PATH=/root/.openclaw/extensions/dingtalk/node_modules
ExecStart=/usr/bin/node /root/yoyoo/backend/scripts/dingtalk_stream_forwarder.cjs
Restart=always
RestartSec=3
# Avoid unbounded log file growth; use journald rotation instead.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
