# Yoyoo/project/backend/deploy/yoyoo-memory-maintenance.timer
[Unit]
Description=Run Yoyoo memory maintenance daily

[Timer]
OnCalendar=*-*-* 04:20:00
Persistent=true
Unit=yoyoo-memory-maintenance.service

[Install]
WantedBy=timers.target

