# Yoyoo/project/backend/deploy/yoyoo-memory-sidecar.service
[Unit]
Description=Yoyoo Memory Sidecar API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/yoyoo/backend
EnvironmentFile=-/root/yoyoo/backend/.env
ExecStart=/root/yoyoo/backend/.venv/bin/uvicorn app.services.memory_sidecar_http:app --host 127.0.0.1 --port 8787
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
