# Yoyoo/project/backend/deploy/yoyoo-memory-maintenance.service
[Unit]
Description=Yoyoo Memory Maintenance
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/yoyoo/backend
EnvironmentFile=-/root/yoyoo/backend/.env
ExecStart=/root/yoyoo/backend/.venv/bin/python scripts/memory_maintenance.py --memory-file /root/yoyoo/backend/data/yoyoo_memory.json --archive-dir /root/yoyoo/backend/data/archive --report-dir /root/yoyoo/backend/data/reports

